# piko-watchdog.service
[Unit]
Description=Piko Hosts Watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/piko-watchdog
